# deployment/systemd/milkbot-rollups.timer
[Unit]
Description=Milkbot Analytics Rollup Timer
Documentation=https://github.com/milkbotai/Kalshi

[Timer]
# Run daily at 1 AM ET (after market close)
OnCalendar=*-*-* 01:00:00
# Run missed jobs on boot
Persistent=true
Unit=milkbot-rollups.service

[Install]
WantedBy=timers.target
